# launchd/com.firstmate.linear-events.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.firstmate.linear-events</string>
  <key>ProgramArguments</key>
  <array>
    <string>__NODE__</string>
    <string>__FIRSTMATE_ROOT__/bin/fm-linear-event-server.mjs</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__FIRSTMATE_ROOT__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>FM_LINEAR_EVENT_HOST</key>
    <string>127.0.0.1</string>
    <key>FM_LINEAR_EVENT_PORT</key>
    <string>4481</string>
    <key>FM_LINEAR_EVENT_STATE</key>
    <string>__FIRSTMATE_ROOT__/state/linear-events</string>
    <key>FM_LINEAR_EVENT_SECRET_FILE</key>
    <string>__FIRSTMATE_ROOT__/state/linear-events/webhook-secret</string>
    <key>FM_LINEAR_EVENT_WORKER</key>
    <string>__FIRSTMATE_ROOT__/bin/fm-linear-event-worker.sh</string>
    <key>FM_LINEAR_EVENT_CLAUDE</key>
    <string>__CLAUDE__</string>
    <key>FM_ROOT_OVERRIDE</key>
    <string>__FIRSTMATE_LIVE_ROOT__</string>
    <key>PATH</key>
    <string>__PATH__</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>__FIRSTMATE_ROOT__/state/linear-events/service.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>__FIRSTMATE_ROOT__/state/linear-events/service.stderr.log</string>
</dict>
</plist>
